# rcon/xinetd.d/rcon
# description: An xinetd internal service which echo's characters back to
# clients.
service rcon
{
	port            = 13998
	disable		= no
#	type		= INTERNAL
#	id		= echo-stream
	socket_type	= stream
	protocol	= tcp
	user		= ryzcom
	wait		= no
	server		= /usr/local/bin/rcon
	server_args	= /home/ryzcom/etc /home/ryzcom/usersdata /home/ryzcom/bin_rcon
}
